# lxqt/lxqt-overlay/etc/skel/.config/autostart/cmst.desktop
[Desktop Entry]
Exec=cmst
Name=cmst
Type=Application
Version=1.0
Icon=cmst
X-LXQt-Need-Tray=true
