# lxqt/lxqt-overlay/etc/skel/.config/autostart/connman-applet.desktop
[Desktop Entry]
Exec=/usr/bin/lxqt-connman-applet
Name=connman-applet
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
