# xfce-openrc/overlay-livecd/etc/local.d/rc.start
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# Source /etc/conf.d/xdm if it exists
[ -e /etc/conf.d/xdm ] && . /etc/conf.d/xdm

# Some hackery for (re)starting lxdm
if [ "$DISPLAYMANAGER" == "lxdm" ]; then
	# Restart xdm as lxdm does not get started correctly the first time
	rc-service xdm restart
fi
