# lxqt-minimal/lxqt-minimal-overlay/etc/skel/.config/autostart/cmst.desktop
[Desktop Entry]
Exec=cmst -m
Name=Cmst
Icon=cmst
Type=Application
Version=1.1
X-LXQt-Need-Tray=true

