# shared/overlay-livecd/etc/skel/.config/autostart/live-session.desktop
[Desktop Entry]
Exec=live-session
Name=desktop_live_items
Type=Application
Version=1.1

